# adkill/files/adkill
#!/bin/sh /etc/rc.common

START=25

start() {
	insmod /lib/modules/3.10.36/adkill.ko
	chmod 777 /etc/adkill/advkill-conf.sh
	/etc/adkill/advkill-conf.sh &
}

